# systemd/.config/systemd/user/sync-mail.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/home/jamie/scripts/sync-mail -a
