# systemd/.config/systemd/user/sync-mail.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnUnitActiveSec=1m
Unit=sync-mail.service

[Install]
WantedBy=timers.target
